# build/linux/rootfs/stage/var/lib/dpkg/info/busybox-mdev.list
/.
/etc
/etc/init.d
/etc/init.d/mdev
/etc/mdev.conf
/etc/mdev
/etc/mdev/mdev-mount.sh
/etc/mdev/find-touchscreen.sh
