# build/linux/rootfs/targetroot/var/lib/dpkg/info/busybox-syslog.list
/.
/etc
/etc/init.d
/etc/init.d/syslog.busybox
/etc/syslog-startup.conf.busybox
/etc/syslog.conf.busybox
